# agent/com.meridian.merdx.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.meridian.merdx</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/zach/Downloads/meridian/agent/_merd-post.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>17</integer></dict>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>17</integer></dict>
    <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>17</integer></dict>
    <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>17</integer></dict>
  </array>
  <key>StandardOutPath</key>
  <string>/Users/zach/Downloads/meridian/agent/_autopilot.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/zach/Downloads/meridian/agent/_autopilot.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
